# 3dprinter_debian/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

# Command added by /etc/cont-init.d/mjpg-streamer.sh.
