# tsd/rootfs/etc/services.d/web/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting web..."

source /envvars.sh

cd /data/web/

sh -c "python manage.py collectstatic --noinput && python manage.py migrate && python manage.py runserver --noreload 0.0.0.0:3334"