# tsd/rootfs/etc/services.d/tasks/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting tasks..."

source /envvars.sh

cd /data/web/

sh -c "celery -A config worker --beat -l info -c 2 -Q realtime,celery"
