# octoprint_dev/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

if [ -f /config/nginx.conf ]; then
    bashio::log.info "Starting nginx"
    nginx -c /config/nginx.conf
fi
