# tsd/rootfs/etc/services.d/tinyproxy/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting reverse proxy for Ingress..."

tinyproxy -d -c /config/tinyproxy.conf