# octoprint-next-dev/rootfs/etc/services.d/caddy/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Terminate Caddy"
