# octoprint-next-dev/rootfs/etc/services.d/octoprint/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Finish up OctoPrint"

bashio::addon.stop