# octoprint-next-dev/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting octoprint..."

octoprint serve --iknowwhatimdoing --basedir /data/config/octoprint --host 0.0.0.0 --port 80
