# octoprint_py3/FS/etc/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 1MB
loglevel = debug
pidfile = /tmp/supervisord.pid
nodaemon = false
user = root
childlogdir = /tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:haproxy]
command = /usr/sbin/haproxy -db -f /etc/haproxy/haproxy.cfg
stdout_logfile_maxbytes = 1MB

[program:dbus-daemon]
command = dbus-daemon --system --nofork
stdout_logfile_maxbytes = 1MB

[program:avahi-daemon]
command = avahi-daemon
stdout_logfile_maxbytes = 1MB

[program:octoprint]
command = octoprint serve --iknowwhatimdoing --basedir %%basedir%% --host 0.0.0.0 --port 80
autorestart = true
stdout_logfile_maxbytes = 1MB

[program:mjpg-streamer]
autostart = true
command = mjpg_streamer -i "%%mjpg_input%%" -o "%%mjpg_output%%"
stdout_logfile_maxbytes = 1MB
