# klipper/rootfs/etc/services.d/fluidd/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Fluidd."

caddy run -config /etc/caddy/Caddyfile
