# ustreamer/rootfs/etc/services.d/tinyproxy-web/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting tinyproxy for web..."

tinyproxy -d -c tinyproxy-web.conf