# pioled/rootfs/etc/services.d/pioled/finish
#!/usr/bin/with-contenv bashio

exec python /code/end.py