# pioled/rootfs/etc/services.d/pioled/run
#!/usr/bin/with-contenv bashio

exec python /code/run.py