# 3dprinter/rootfs/etc/services.d/proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start proxy service
# ==============================================================================

bashio::log.info "Starting proxy..."

# exec caddy run -config /etc/caddy/Caddyfile
exec nginx