# octoprint-base/rootfs/etc/services.d/proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start proxy service
# ==============================================================================

bashio::log.info "Starting proxy..."

exec haproxy -db -f /etc/haproxy/haproxy.cfg
