# octoprint-base/rootfs/etc/services.d/proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start proxy service
# ==============================================================================

bashio::log.info "Starting proxy..."

exec nginx -g "daemon off; error_log /dev/stdout debug;"
