# klipper/rootfs/etc/services.d/klipper/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Klipper."

KLIPPER_PYTHON=/root/klippy-env/bin/python
KLIPPY=/root/klipper/klippy/klippy.py
KLIPPER_CFG=${HOME}/printer.cfg
KLIPPY_LOGFILE=/tmp/klippy.log

touch ${KLIPPY_LOGFILE}

${KLIPPER_PYTHON} ${KLIPPY} ${KLIPPER_CFG} -l ${KLIPPY_LOGFILE}
