# octoprint-next-dev/rootfs/etc/services.d/proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start proxy service
# ==============================================================================

bashio::log.info "Starting proxy..."

caddy run -config /etc/caddy/Caddyfile
# exec haproxy -p /tmp/haproxy.pid -f /etc/haproxy/haproxy.cfg -db
# exec tinyproxy -d -c /etc/tinyproxy/tinyproxy-ingress.conf
