# 3dprinter/rootfs/etc/services.d/klipper/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Finish up Klipper..."
