# 3dprinter/rootfs/etc/services.d/klipper/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Klipper..."

KLIPPER_PYTHON=/data/python/Klipper/bin/python
KLIPPY=/data/klipper/klippy/klippy.py
KLIPPER_CFG=/data/config/octoprint/klipper/printer.cfg
KLIPPY_LOGFILE=/tmp/klippy.log
KLIPPY_UDS=/tmp/klippy_uds

touch ${KLIPPY_LOGFILE}
touch ${KLIPPER_CFG}

exec ${KLIPPER_PYTHON} ${KLIPPY} ${KLIPPER_CFG} -l ${KLIPPY_LOGFILE} -a ${KLIPPY_UDS}
