# octoprint_dev/rootfs/etc/services.d/haproxy/run
#!/usr/bin/with-contenv bashio

if [ -f /config/haproxy.cfg ]; then
    bashio::log.info "Starting HAProxy"
    haproxy -db -f /config/haproxy.cfg
fi
