# tsd/rootfs/etc/services.d/ingress/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting proxy for Ingress..."

/usr/sbin/haproxy -db -f /etc/haproxy/haproxy.cfg