# octoprint_dev/rootfs/etc/services.d/tinyproxy-ingress/run
#!/usr/bin/with-contenv bashio

if [ -f /config/tinyproxy-ingress.conf ]; then
    bashio::log.info "Starting Tinyproxy for Ingress..."
    tinyproxy -d -c /config/tinyproxy-ingress.conf
fi
