# octoprint_next/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio

## Preparations of OctoPrint is done in /etc/cont-init.d/octoprint.sh

# Launch OctoPrint
sleep 5 && octoprint serve --iknowwhatimdoing --basedir /data/octoprint --host 0.0.0.0 --port 80
