# octoprint_dev/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio

## Preparations of OctoPrint is done in /etc/cont-init.d/octoprint.sh

## Read values from configuration.
OCTOPRINT_BASEDIR="/config/octoprint"$(bashio::config 'config_folder_suffix')

octoprint serve --iknowwhatimdoing --basedir "${OCTOPRINT_BASEDIR}" --host 0.0.0.0 --port 80
