# octoprint-extended/rootfs/etc/services.d/proxy-ingress/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start proxy service
# ==============================================================================

bashio::log.info "Starting proxy (Ingress)..."

exec tinyproxy -d -c /data/proxy/ingress.conf
