# tsd/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting reverse proxy for Ingress..."

nginx -c /config/nginx.conf
