# 3dprinter/rootfs/etc/services.d/mjpg-streamer/finish
#!/usr/bin/with-contenv bashio

bashio::log.notice "Finish up mjpg-streamer"
