# 3dprinter/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

# Command is added by /etc/cont-init.d/mjpg-streamer.sh.
