# octoprint_dev/rootfs/etc/services.d/tinyproxy-webui/run
#!/usr/bin/with-contenv bashio

if [ -f "/config/octoprint"$(bashio::config 'config_folder_suffix')/tinyproxy-webui.conf ]; then
    bashio::log.info "Starting Tinyproxy for Ingress..."
    tinyproxy -d -c "/config/octoprint"$(bashio::config 'config_folder_suffix')/tinyproxy-webui.conf
fi
