# octoprint_dev/rootfs/etc/services.d/tinyproxy-webui/run
#!/usr/bin/with-contenv bashio

if [ -f /config/tinyproxy-webui.conf ]; then
    bashio::log.info "Starting Tinyproxy for Ingress..."
    tinyproxy -d -c /config/tinyproxy-webui.conf
fi
