# mjpg-streamer/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: OctoPrint (Dev)
# Starts server to serve the app
# ==============================================================================

bashio::log.info $(bashio::addon.name)
bashio::log.info $(bashio::addon.version)
bashio::log.info $(bashio::addon.hostname)

mjpg_streamer -i "$(bashio::config 'mjpg_input')" -o "output_http.so -w /www_mjpg -p 8000"
