# octoprint-next-dev/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start nginx service
# ==============================================================================

bashio::log.info "Starting nginx..."

exec nginx
