# OctoPrint-base/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Caddy service
# ==============================================================================

bashio::log.info "Starting Caddy..."

caddy run -config /etc/caddy/Caddyfile
