# octoprint_dev/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Caddy ..."

ingress_entry=$(bashio::addon.ingress_entry)

# TODO: This is useful for debugging, may not be suitable for final version.
if [ -f "/config/octoprint"$(bashio::config 'config_folder_suffix')/Caddyfile ]; then
    sed -i "s#%%base_path%%#${ingress_entry}#g" "/config/octoprint"$(bashio::config 'config_folder_suffix')/Caddyfile
    caddy run --watch -config "/config/octoprint"$(bashio::config 'config_folder_suffix')/Caddyfile
else
    sed -i "s#%%base_path%%#${ingress_entry}#g" /etc/services.d/caddy/Caddyfile
    caddy run -config /etc/services.d/caddy/Caddyfile
fi
