# octoprint_dev/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio

bashio::log.info "caddy"

caddy run -config /config/Caddyfile
# TODO: Use Caddyfile in this folder, not from config.
# caddy run -config Caddyfile
