# octoprint_dev/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio

if [ -f /config/Caddyfile ]; then
    bashio::log.info "Starting Caddy ..."
    caddy run -config /config/Caddyfile
fi
