# octoprint_dev/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio

if [ -f "/config/octoprint"$(bashio::config 'config_folder_suffix')/Caddyfile ]; then
    bashio::log.info "Starting Caddy ..."
    caddy run -config "/config/octoprint"$(bashio::config 'config_folder_suffix')/Caddyfile
    # caddy run -config /etc/services.d/caddy/Caddyfile
fi
