# 3dprinter/rootfs/etc/services.d/moonraker/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Moonraker..."

MOONRAKER_PYTHON=/data/python/moonraker/bin/python
MOONRAKER=/data/src/moonraker/moonraker/moonraker.py
MOONRAKER_CFG=/data/config/klipper/moonraker.conf
MOONRAKER_LOGFILE=/tmp/moonraker.log

touch ${MOONRAKER_LOGFILE}
touch ${MOONRAKER_CFG}

exec ${MOONRAKER_PYTHON} ${MOONRAKER} -c ${MOONRAKER_CFG} -l ${MOONRAKER_LOGFILE}
