# octoprint-next-dev/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

# join with separator
join_ws() { local d=$1 s=$2; shift 2 && printf %s "$s${@/#/$d}"; }

mjpg_input=$(bashio::config 'mjpg_input')
mjpg_output=$(bashio::config 'mjpg_output')

exec mjpg_streamer ${join_ws '-o ' mjpg_output} ${join_ws '-i ' mjpg_input}
