# octoprint-next-dev/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

# This should be modified to allow multiple input terms from bashio

delim=" -i "
printf -v var "$delim \"%s\"" "${bashio::config 'mjpg_input'[@]}"
var="${var%$delim}"

echo $var

mjpg_streamer -o "output_http.so -w /www_mjpg -p 8080" $var
