# octoprint-next-dev/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

exec mjpg_streamer %%mjpg_arguments%%
