# octoprint-next-dev/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

mjpg_input=$(bashio::config 'mjpg_input')
m_input=$(printf " -i \"%s\"" "${mjpg_input[@]}")
echo $m_input
mjpg_output=$(bashio::config 'mjpg_output')
m_output=$(printf " -o \"%s\"" "${mjpg_output[@]}")
echo $m_output

exec mjpg_streamer $m_output $m_input
