# octoprint-next-dev/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

mjpg_arguments=$(bashio::config 'mjpg_arguments')
echo $mjpg_arguments
mjpg_arguments=${mjpg_arguments[*]}
echo $mjpg_arguments
mjpg_arguments=${mjpg_arguments// / }
echo $mjpg_arguments
exec mjpg_streamer $mjpg_arguments
