# octoprint-next-dev/rootfs/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mjpg-streamer..."

# echo $(printf ' -i \"%s\"' "${bashio::config 'mjpg_input'[@]}")
# exec mjpg_streamer -o "output_http.so -w /www_mjpg -p 8080" $(printf ' -i \"%s\"' "${bashio::config 'mjpg_input'[@]}")
exec mjpg_streamer -o "output_http.so -w /www_mjpg -p 8080" ${mjpg_input}[@]
