# 3dprinter/rootfs/etc/services.d/klipper-linux-host/finish
#!/usr/bin/with-contenv bashio

bashio::log.notice "Finish up Klipper MCU Linux..."
